# scripts/carbon-lens-sync.service
[Unit]
Description=Carbon Lens MongoDB to PostgreSQL Sync
After=network.target postgresql.service mongod.service

[Service]
Type=oneshot
User=aneeq
WorkingDirectory=/home/aneeq/Documents/Fiver/react-dashboard/CL_site/mongo-pg
EnvironmentFile=/home/aneeq/Documents/Fiver/react-dashboard/CL_site/mongo-pg/.env
ExecStart=/home/aneeq/Documents/Fiver/react-dashboard/CL_site/mongo-pg/.venv/bin/meltano run tap-mongodb target-postgres
StandardOutput=append:/home/aneeq/Documents/Fiver/react-dashboard/CL_site/mongo-pg/logs/sync.log
StandardError=append:/home/aneeq/Documents/Fiver/react-dashboard/CL_site/mongo-pg/logs/sync.log

[Install]
WantedBy=multi-user.target
