# scripts/carbon-lens-sync.timer
[Unit]
Description=Run Carbon Lens sync every hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
